# deploy/systemd/daily-scheduled-worker.timer
[Unit]
Description=Run the Daily scheduled summary worker once per minute
Documentation=https://github.com/poppyseedcake/daily

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
Unit=daily-scheduled-worker.service

[Install]
WantedBy=timers.target

