# deploy/systemd/daily-backup.service
[Unit]
Description=Daily verified SQLite backup
Documentation=https://github.com/poppyseedcake/daily

[Service]
Type=oneshot
User=daily
Group=daily
WorkingDirectory=/srv/daily/current
Environment=NODE_ENV=production
EnvironmentFile=/etc/daily/daily.env
Environment=DAILY_SYSTEMD_BACKUP_UNIT=true
ExecStartPre=/usr/bin/node scripts/validate-production-environment.mjs
ExecStart=/usr/bin/node build/worker/runSqliteBackupCommand.js daily
StandardOutput=journal
StandardError=journal
SyslogIdentifier=daily-backup
UMask=0027

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX
SystemCallArchitectures=native
ReadWritePaths=/var/lib/daily /var/backups/daily
