# deploy/systemd/daily-scheduled-worker.service
[Unit]
Description=Daily scheduled summary worker
Documentation=https://github.com/poppyseedcake/daily
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=daily
Group=daily
WorkingDirectory=/srv/daily/current
Environment=NODE_ENV=production
EnvironmentFile=/etc/daily/daily.env
ExecStartPre=/usr/bin/node scripts/validate-production-environment.mjs
ExecStart=/usr/bin/npm run worker:scheduled-delivery
StandardOutput=journal
StandardError=journal
SyslogIdentifier=daily-scheduled-worker
UMask=0027

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
ReadWritePaths=/var/lib/daily

